# deploy/rental-system.service
[Unit]
Description=Rental System API
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/rental-system/current
EnvironmentFile=/var/www/rental-system/shared/.env
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5
User=www-data
Group=www-data
StandardOutput=append:/var/log/rental-system/app.log
StandardError=append:/var/log/rental-system/error.log

[Install]
WantedBy=multi-user.target
